# root/etc/services.d/filebrowser/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# ==============================================================================
# File Browser Service Definition
# Optimized for: Plex Media Server Ecosystem
# Role: Senior AppSec & Full-Stack Engineering Standard
# Security Level: Hardened + Regression Mitigated + UX Enhanced
# ==============================================================================

# --- Configuration & Constants ---
readonly FB_BINARY="/usr/local/bin/filebrowser"
readonly USER="plex"
readonly GROUP="plex"
readonly DATA_DIR="/data"
readonly CONFIG_DIR="/config"
readonly DB_PATH="${CONFIG_DIR}/filebrowser.db"
# [REGRESSION FIX] Parameterized Port (Default: 9090)
readonly PORT="${FB_PORT:-9090}"
readonly BIND_ADDRESS="0.0.0.0"
# [REGRESSION FIX] Initialization Lock
readonly INIT_LOCK="${CONFIG_DIR}/.fb_init.lock"
# [UX] Branding Directory
readonly BRANDING_DIR="${CONFIG_DIR}/branding"

# --- Credentials (Secrets Management) ---
# Priority:
# 1. Docker Secret File (FB_ADMIN_PASSWORD_FILE)
# 2. Environment Variable (FB_ADMIN_PASSWORD)
# 3. Default (admin)

ADMIN_USER="${FB_ADMIN_USER:-admin}"
ADMIN_PASS_DEFAULT="admin"

# Resolve Password
if [ -n "$FB_ADMIN_PASSWORD_FILE" ] && [ -f "$FB_ADMIN_PASSWORD_FILE" ]; then
    ADMIN_PASS=$(cat "$FB_ADMIN_PASSWORD_FILE")
elif [ -n "$FB_ADMIN_PASSWORD" ]; then
    ADMIN_PASS="$FB_ADMIN_PASSWORD"
else
    ADMIN_PASS="$ADMIN_PASS_DEFAULT"
fi

# Sanity check
if [ -z "$ADMIN_PASS" ]; then
    echo "[CRITICAL] Password cannot be empty. Reverting to default."
    ADMIN_PASS="$ADMIN_PASS_DEFAULT"
fi

# --- Helper Functions ---

log_info() {
    echo "[INFO] $(date '+%Y-%m-%d %H:%M:%S') - $1"
}

log_warn() {
    echo "[WARN] $(date '+%Y-%m-%d %H:%M:%S') - $1"
}

log_error() {
    echo "[ERROR] $(date '+%Y-%m-%d %H:%M:%S') - $1"
}

ensure_permissions() {
    log_info "Verificando permisos (Principio de Mínimo Privilegio)..."
    
    if [ -d "$DATA_DIR" ]; then
        current_user=$(stat -c '%U' "$DATA_DIR" 2>/dev/null || echo "unknown")
        if [ "$current_user" != "$USER" ]; then
             chown "$USER:$GROUP" "$DATA_DIR" || log_warn "No se pudo cambiar el dueño de $DATA_DIR"
        fi
    fi

    if [ -d "$CONFIG_DIR" ]; then
        chown -R "$USER:$GROUP" "$CONFIG_DIR"
    fi
    
    if [ -f "$DB_PATH" ]; then
        chown "$USER:$GROUP" "$DB_PATH"
    fi
}

# [UX] Ensure Plex-Themed Branding
ensure_branding() {
    log_info "Verificando identidad visual (UX)..."
    
    if [ ! -d "$BRANDING_DIR" ]; then
        mkdir -p "$BRANDING_DIR"
        chown "$USER:$GROUP" "$BRANDING_DIR"
    fi

    # Plex Dark Theme CSS
    # Colors: Background #1f1f1f, Sidebar #18191c, Accent #e5a00d
    CSS_FILE="${BRANDING_DIR}/custom.css"
    
    cat > "$CSS_FILE" <<EOF
/* Plex Integration Theme - Generated by Antigravity */
:root {
  --text-main: #e5e5e5;
  --text-primary: #e5e5e5;
  --text-secondary: #999;
  --bg-main: #1f1f1f;
  --bg-secondary: #282a2d;
  --action-main: #e5a00d; /* Plex Orange */
  --action-secondary: #cc8e00;
  --selection: rgba(229, 160, 13, 0.2);
  --border: #3a3d42;
}
body { background-color: var(--bg-main) !important; color: var(--text-main) !important; }
header { background-color: var(--bg-secondary) !important; border-bottom: 1px solid var(--border) !important; }
.sidebar { background-color: #18191c !important; border-right: 1px solid var(--border) !important; }
input, select, textarea { background-color: #323539 !important; border: 1px solid var(--border) !important; color: white !important; }
button.action { background-color: var(--action-main) !important; color: black !important; font-weight: bold; }
button.action:hover { background-color: var(--action-secondary) !important; }
/* Hide generic footer to reduce visual noise */
.credits { display: none !important; }
EOF

    # Ensure permissions for the new css file
    chown "$USER:$GROUP" "$CSS_FILE"
}

initialize_database() {
    # [REGRESSION FIX] Lock Mechanism
    if [ -f "$INIT_LOCK" ]; then
        if [ "$(find "$INIT_LOCK" -mmin +1 2>/dev/null)" ]; then
            log_warn "Eliminando lock file obsoleto (>1 min)."
            rm -f "$INIT_LOCK"
        else
            log_warn "Bloqueo de inicialización detectado ($INIT_LOCK). Saltando init."
            return 0
        fi
    fi

    touch "$INIT_LOCK"
    # shellcheck disable=SC2064
    trap "rm -f '$INIT_LOCK'" RETURN

    if [ ! -f "$DB_PATH" ]; then
        log_info "Inicializando base de datos nueva..."
        # Init DB only
        if ! "$FB_BINARY" config init --database "$DB_PATH" > /dev/null 2>&1; then
            log_error "Fallo al inicializar DB."
            exit 1
        fi
    else
        log_info "Base de datos existente detectada."
    fi

    # [FIX] Always ensure Admin User exists and has correct password (Idempotency)
    log_info "Asegurando usuario administrador ($ADMIN_USER)..."
    
    # Intenta actualizar contraseña primero (si el usuario ya existe)
    # Capturamos output para debug
    if update_out=$("$FB_BINARY" users update "$ADMIN_USER" --password "$ADMIN_PASS" --perm.admin --database "$DB_PATH" 2>&1); then
        log_info "Usuario $ADMIN_USER actualizado correctamente."
    else
        # Si falla (ej: usuario no existe), intentamos crear
        # Regex extendida para capturar "no rows" (sqlite) O "does not exist" (app logic)
        if echo "$update_out" | grep -qE "no rows|does not exist"; then
             log_info "El usuario no existe (update falló). Intentando crear nuevo..."
             if add_out=$("$FB_BINARY" users add "$ADMIN_USER" "$ADMIN_PASS" --perm.admin --database "$DB_PATH" 2>&1); then
                 log_info "Usuario $ADMIN_USER creado exitosamente."
             else
                 log_error "Fallo al crear usuario: $add_out"
             fi
        else
             # Hubo otro error en el update, intentamos 'add' por si acaso (fallback brutal)
             log_warn "Fallo update desconocido ($update_out). Intentando creación forzada..."
             if add_fallback_out=$("$FB_BINARY" users add "$ADMIN_USER" "$ADMIN_PASS" --perm.admin --database "$DB_PATH" 2>&1); then
                  log_info "Usuario asegurado mediante fallback."
             else
                  log_error "NO SE PUDO crear el usuario administrador. Error: $add_fallback_out"
             fi
        fi
    fi

    if [ "$ADMIN_PASS" == "admin" ]; then
         log_warn "SECURITY ALERT: Credenciales por defecto en uso (admin/admin)."
    fi
}

configure_application() {
    log_info "Aplicando configuración blindada..."
    
    "$FB_BINARY" config set --database "$DB_PATH" \
        --address "$BIND_ADDRESS" \
        --port "$PORT" \
        --root "$DATA_DIR" \
        --auth.method=json \
        --signup=false \
        --branding.name "Plex File Manager" \
        --branding.disableExternal \
        --branding.files "$BRANDING_DIR" \
        > /dev/null 2>&1

    # [REGRESSION FIX] Robust Version Compatibility Check
    if ! err_msg=$("$FB_BINARY" config set --database "$DB_PATH" --auth.min-password-length 10 2>&1); then
         if ! "$FB_BINARY" config set --database "$DB_PATH" --minimumPasswordLength 10 > /dev/null 2>&1; then
             log_warn "No se pudo establecer política de contraseñas. Error: $err_msg"
         fi
    fi
    
    chown "$USER:$GROUP" "$DB_PATH"
}

# --- Main Execution Flow ---

cleanup() {
    unset ADMIN_PASS
    if [ -f "$INIT_LOCK" ]; then rm -f "$INIT_LOCK"; fi
}
trap cleanup EXIT

log_info "Iniciando servicio File Browser (Secured + Mitigated + Themed)..."

if [ ! -f "$FB_BINARY" ]; then
    log_error "Binario no encontrado."
    exit 1
fi

ensure_permissions
ensure_branding
initialize_database
configure_application

log_info "Ejecutando proceso principal en puerto $PORT..."

unset ADMIN_PASS
unset FB_ADMIN_PASSWORD

exec s6-setuidgid "$USER" "$FB_BINARY" --database "$DB_PATH"